# supervisor_template.conf
[program:PROGRAMNAME]
command=COMMAND
user=USER
autostart=true
autorestart=true
stderr_logfile=/var/log/PROGRAMNAME.access.log
stdout_logfile=/var/log/PROGRAMNAME.error.log